# modules/desktop/nixos/kde/autostart/vesktop.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=Element
Exec=vesktop
StartupNotify=false
Terminal=false
