# modules/desktop/kde/autostart/element-desktop.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=Element
Exec=element-desktop --hidden
StartupNotify=false
Terminal=false
